# topic.conf
[program:createTopic]
command=/usr/bin/topic.sh
autostart=true
autorestart=false